# mavproxy.service
[Unit]
Description=MAVProxy telemetry bridge
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/run/mavproxy
RuntimeDirectory=mavproxy
ExecStart=/usr/bin/env mavproxy.py --master=udp:127.0.0.1:14550 --out=/dev/ttyUSB0 --out=udp:0.0.0.0:14551 --out=udp:0.0.0.0:14552 --state-basedir=/run/mavproxy
Restart=on-failure
RestartSec=2
KillSignal=SIGINT
StandardOutput=journal
StandardError=inherit

[Install]
WantedBy=multi-user.target
